# dot_config/systemd/user/tasknotify-check.service
[Unit]
Description=tasknotify — task notification checker
# Do not use PartOf=graphical-session.target here; that would prevent the
# service from running in standalone WM setups (i3, bspwm, sway, etc.) that
# don't fully activate graphical-session.target via a display manager.
# The script itself recovers DBUS_SESSION_BUS_ADDRESS and DISPLAY at runtime.

[Service]
Type=oneshot
ExecStart=/home/phong/.local/bin/tasknotify.sh check

# Provide a known-good DBUS path; the script will override if needed.
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

StandardOutput=append:/home/phong/.local/share/tasknotify/tasknotify.log
StandardError=append:/home/phong/.local/share/tasknotify/tasknotify.log

[Install]
WantedBy=default.target
